# freebsd.in
#!/bin/sh
#
# PROVIDE: notify
# REQUIRE: NETWORKING
# KEYWORD:

. /etc/rc.subr

name="notify"
rcvar="notify_enable"

load_rc_config $name

: ${notify_enable:="NO"}
: ${notify_username:="%%USERNAME%%"}
: ${notify_logfile:="/var/log/notify.log"}

notify_command="/usr/local/bin/notify"
pidfile="/var/run/notify.pid"
command="/usr/sbin/daemon"
command_args="-o ${notify_logfile} -P ${pidfile} -r -u ${notify_username} ${notify_command} --daemon"
required_files="/usr/local/etc/notify.token"
start_precmd="mkdir -p /var/run/notify && chown ${notify_username} /var/run/notify && chmod 755 /var/run/notify"
stop_postcmd="rm -f /var/run/notify/sock"

run_rc_command "$1"
